# config/systemd/user/mcp-hub.service
[Unit]
Description=MCP Hub (aggregator for MCP servers)
After=default.target

[Service]
Type=simple
# Ensure 1Password CLI and mise shims are on PATH so ${cmd: op read ...} resolves.
Environment=PATH=%h/.bin:%h/.local/share/mise/shims:%h/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=%h/.local/share/mise/shims/mcp-hub --port 37373 --config %h/.config/mcphub/servers.json
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
