# config/systemd/user/mcp-hub.service
[Unit]
Description=MCP Hub (aggregator for MCP servers)
After=default.target

[Service]
Type=simple
# Wrapper sets PATH and sources OP_SERVICE_ACCOUNT_TOKEN; same wrapper is used
# by the macOS LaunchAgent at config/launchd/com.user.mcp-hub.plist.
ExecStart=%h/.bin/mcp-hub-run
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
